# deploy/agentzero.service
[Unit]
Description=AgentZero Telegram Bot
After=network.target

[Service]
Type=simple
User=agentzero
Group=agentzero
WorkingDirectory=/opt/agentzero
EnvironmentFile=/opt/agentzero/.env
ExecStart=/opt/agentzero/venv/bin/uvicorn agentzero.main:app --host 127.0.0.1 --port 8080 --log-level info
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
