# deploy/agentzero.service
[Unit]
Description=AgentZero Telegram Bot
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/production/AgentZero
EnvironmentFile=/var/www/production/AgentZero/.env
ExecStart=/var/www/production/AgentZero/venv/bin/uvicorn agentzero.main:app --host 127.0.0.1 --port 8080 --log-level info
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
